# deploy/customer-service.service
[Unit]
Description=Customer Service (Next.js)
After=network.target

[Service]
User=portal-system
Group=portal-system
WorkingDirectory=/opt/customer-service
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run start -- --hostname 127.0.0.1 --port 3001
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
